# 10_useless_lab_site/deploy/useless-lab.service
[Unit]
Description=Useless Lab static mini games site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Raspberry_Pi/10_useless_lab_site
ExecStart=/usr/bin/python3 -m http.server 8011 --bind 0.0.0.0 --directory web
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
