# 02_autorun_monitor_recovery/dashboard.service
[Unit]
Description=Raspberry Pi FastAPI Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Raspberry_Pi/01_web_dashboard_fastapi
EnvironmentFile=/home/user/Raspberry_Pi/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/user/Raspberry_Pi/01_web_dashboard_fastapi/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
